# deploy/spawner.service
[Unit]
Description=claude_spawner voice server
Documentation=https://github.com/bam/claude_spawner
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=bam
Group=bam
WorkingDirectory=/data/claude_spawner/server
EnvironmentFile=/home/bam/.config/claude_spawner/spawner.env
Environment=HOME=/home/bam
Environment=PATH=/home/bam/.local/bin:/usr/local/bin:/usr/bin:/bin

# Rebuild on every (re)start so it runs current code, then exec the binary so
# systemd supervises the real process (not `go run`'s child).
ExecStartPre=/usr/bin/go build -o /data/claude_spawner/server/spawner .
ExecStart=/data/claude_spawner/server/spawner

Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
